# oraQuadra_Nano_Ver_1.5/oraQuadra_Nano_Ver_1.5/LEGGI_ISTRUZIONI/RADAR_CLIENT_ILI9341/radar_monitor_pi.service
[Unit]
Description=Radar Monitor Controller per Raspberry Pi
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/radar_monitor_pi.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
